# Headless/vagrant/run-radar.cron
* * * * * vagrant /vagrant/run-radar.sh
